# deploy/web_admin/linkhive-admin.service
[Unit]
Description=LinkHive Web Admin
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/linkhive.conf
ExecStart=/opt/linkhive/venv/bin/python /usr/local/bin/linkhive_admin.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
