# deploy/sms_forwarder/sms-forwarder.service
[Unit]
Description=Forward incoming SMS messages via Apprise
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/linkhive/venv/bin/python /usr/local/bin/sms_forwarder.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
